# images/stackstorm/conf/local.conf
description "Local startup"
author "StackStorm"

start on started st2api
stop on shutdown

script

logger Starting local init...

st2ctl reload --register-all

logger Ending local init

end script
