# images/stackstorm-xenial/supervisord.d/sshd.conf

[program:sshd]
command = /st2-docker/bin/run-sshd.sh
startsecs = 0
exitcodes = 0
redirect_stderr = true
