# images/stackstorm-xenial/supervisord.d/st2chatops.conf
[program:st2chatops]
command = /st2-docker/bin/run-st2/%(program_name)s.sh
user = st2
redirect_stderr = true
